# etc/mcp-server.service
[Unit]
Description=MCP Server - Model Context Protocol Server
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/www/profullstack.com/mcp-server
ExecStart=/bin/zsh /home/ubuntu/www/profullstack.com/mcp-server/bin/start.sh
Restart=on-failure
RestartSec=5
StandardOutput=append:/var/log/mcp-server.log
StandardError=append:/var/log/mcp-server-error.log
EnvironmentFile=/home/ubuntu/www/profullstack.com/mcp-server/.env
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target